# systemd/priapos.service
[Unit]
Description=Priap.OS - Battle Royale Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pln/Work/Perso/Priap.OS
EnvironmentFile=/home/pln/Work/Perso/Priap.OS/.env
ExecStart=/home/pln/Work/Perso/Priap.OS/.venv/bin/python scripts/br_daemon.py
Restart=on-failure
RestartSec=60

# Graceful shutdown
KillSignal=SIGTERM
TimeoutStopSec=30

# Let daemon handle file logging; journalctl captures stdout/stderr
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
